# Minecraft/crontab
00 00 * * * /home/minecraft/lavahub/restart.sh
